# apps/bob-ws-gateway/bob-ws-gateway.service
[Unit]
Description=Bob WS Gateway (slim relay)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=bob
Group=bob
WorkingDirectory=/opt/bob/ws-gateway
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=GATEWAY_PORT=3002
Environment=BOB_AUTH_BYPASS=true
Environment=BOB_AUTH_BYPASS_USER_ID=default-user
# BOB_AUTH_BYPASS_TOKEN is loaded from /opt/bob/ws-gateway/.env.
EnvironmentFile=/opt/bob/ws-gateway/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bob/ws-gateway
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bob-ws-gateway

[Install]
WantedBy=multi-user.target
