# apps/bob-ws-gateway/bob-ws-gateway.service
[Unit]
Description=Bob WS Gateway (slim relay)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=bob
Group=bob
WorkingDirectory=/opt/bob/ws-gateway
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
Environment=GATEWAY_PORT=3002
# Auth is REAL: browser Better Auth sessions + hashed revocable API keys.
# BOB_AUTH_BYPASS must never appear here — the gateway refuses to boot with
# it set in production (assertNoAuthBypassInProduction). Callers that used
# the bypass form get an API key provisioned via the settings router instead.
# NUDGE_SHARED_SECRET (legacy internal-endpoint ramp) + DATABASE_URL live in
# /opt/bob/ws-gateway/.env; set BOB_ALLOW_LEGACY_NUDGE_SECRET=false there
# once every internal caller carries an API key.
EnvironmentFile=/opt/bob/ws-gateway/.env

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/bob/ws-gateway
PrivateTmp=true

# Resource limits
LimitNOFILE=65536
MemoryMax=512M

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bob-ws-gateway

[Install]
WantedBy=multi-user.target
