# apps/bob-task-runner/bob-task-runner.service
# Reference systemd unit for the Bob Task Runner on hetzner-bob.
# Install: copy to /etc/systemd/system/bob-task-runner.service, then
#   systemctl daemon-reload && systemctl enable --now bob-task-runner
#
# Secrets/config live in the EnvironmentFile (NOT in the repo). At minimum set:
#   BOB_API_URL, BOB_API_KEY, BOB_WORKSPACE_ID  (run reporting)
#   PULSE_API_KEY                                (codex/bizpulse, optional)
# Linear key is read from /home/bob/.linear-key by default.

[Unit]
Description=Bob Task Runner (Linear issue -> codex agent)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bob
WorkingDirectory=/opt/bob/task-runner
EnvironmentFile=/opt/bob/task-runner/.env
Environment=NODE_ENV=production
Environment=PATH=/home/bob/.local/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/node task-runner.js
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
